# recipes-bsp/init/target-system-init/target-system-init.sh.in
#! /bin/sh
### BEGIN INIT INFO
# Provides: Resize root filesystem to fit available disk space
# Description: Start resize helper, then remove own init from runlevels,
# as it is required only for first boot. run settle to provide partition
# table details.
### END INIT INFO

DESC="target system init"

case $1 in
	start)
		echo "Starting $DESC"
		@sbindir@/target-system-init

		echo "remove $DESC from runlevels starting"
		rm @sysconfdir@/rc3.d/S22target-system-init.sh &> /dev/null
		rm @sysconfdir@/rc5.d/S22target-system-init.sh &> /dev/null
		echo "remove $DESC from runlevels completed"
	;;
  *)
		echo "Usage: @sysconfdir@/init.d/target-system-init.sh {start}" >&2
		exit 1
	;;
esac

exit 0

# vim:noet
